# my-ftp-browser/rootfs/etc/s6-overlay/s6-rc.d/ftpapi/finish
#!/command/with-contenv bash

echo "FTP Browser API stopped"
