# my-ftp-browser/rootfs/etc/s6-overlay/s6-rc.d/ftpapi/run
#!/command/with-contenv bash

echo "Starting FTP Browser API..."
cd /app/ftpbrowser/api
exec python3 server.py
