# my-ftp-browser/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv sh
nginx -s stop
