# my-ftp-browser/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bashio
bashio::log.info "Démarrage du serveur web NGINX..."
exec nginx -g "daemon off;"

