# my-ftp-browser/rootfs/etc/services.d/nginx/run
#!/command/with-contenv bashio
bashio::log.info "Démarrage du service NGINX..."

# Vérification de la configuration nginx
nginx -t

# Démarrer nginx en mode premier plan
exec nginx -g "daemon off;"


