# my-ftp-browser/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
set -e

bashio::log.info "Démarrage du serveur web NGINX..."
exec nginx -g "daemon off;"
