# rootfs/etc/services.d/ftp-server/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Démarrage du serveur API FTP..."

# Démarrer le serveur API Python
cd /usr/share/ftpbrowser/api
exec python3 server.py
