# my-ftp-browser/rootfs/etc/services.d/ftp-server/finish
#!/usr/bin/with-contenv sh
pkill -f "python3 /usr/share/ftpbrowser/api/server.py"

