# my-ftp-browser/rootfs/etc/services.d/ftp-server/run
#!/command/with-contenv bashio
bashio::log.info "Démarrage du service API FTP Browser..."

# Configuration Python
export PYTHONUNBUFFERED=1
export PYTHONPATH=/usr/share/ftpbrowser/api

# Attendre que la base de données soit prête
bashio::log.info "Vérification des dépendances..."
timeout 30 bash -c "until python3 -c 'import socket; socket.socket(socket.AF_INET, socket.SOCK_STREAM).connect((\"127.0.0.1\", 5000))' 2>/dev/null; do sleep 1; done" || {
    bashio::log.error "Le service API n'a pas démarré correctement"
    exit 1
}

# Démarrer en surveillant les crashes
while true; do
    bashio::log.info "Lancement du serveur FTP..."
    python3 /usr/share/ftpbrowser/api/server.py \
        >>/proc/1/fd/1 2>>/proc/1/fd/2
    status=$?
    bashio::log.warning "Le serveur a crashé avec le code $status. Redémarrage dans 5s..."
    sleep 5
done
