# my-ftp-browser/rootfs/etc/services.d/ftp-server/run
#!/usr/bin/with-contenv sh
exec python3 /usr/share/ftpbrowser/api/server.py
