# my-ftp-browser/rootfs/etc/services.d/ftp-server/run
#!/usr/bin/env bash
set -e

bashio::log.info "Démarrage du serveur API FTP Browser..."
cd /usr/share/ftpbrowser/api
exec python3 server.py
