# my-ftp-browser/rootfs/etc/services.d/ftp-server/run
#!/command/with-contenv bashio
bashio::log.info "Démarrage du serveur API FTP..."
cd /usr/share/ftpbrowser/api
exec python3 server.py

